# com.claudecrawler.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudecrawler.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kjonekong/pyStcratch/scripts/auto_crawl_and_sync.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/kjonekong/pyStcratch/logs/scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kjonekong/pyStcratch/logs/scheduler_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/kjonekong/pyStcratch</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
